# .docker/cron/crontab
*/10 * * * * root /bin/bash /cron/active.sh >> /logs/active.log 2>&1
